# crontab
*/10 * * * * root ./CronScript.sh >> /var/log/Benchmark-dataUpdation.log 2>&1
# one line at the end should be there
